# vu.kmx.io/etc/rc.d/seuldanslenoir_staging
#!/bin/ksh

daemon="/bin/sh /home/seuldanslenoir-staging/start_"
daemon_user=seuldanslenoir-staging
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
